# lsfusion
##!/bin/sh
# chkconfig: 2345 99 50
# description: lsFusion management script
#
# Adapt the following lines to your configuration:

LSFUSION_HOME=/usr/bin/lsfusion
JAVA_HOME=/usr/lib/jvm/java-8-openjdk-amd64/jre
PID_FILE=$LSFUSION_HOME/jsvc-lsfusion.pid

LSFUSION_OPTS="-Xms512m -Xmx1048m -Xss2m"

CLASSPATH=.:lsfusion-server-1.2.4.jar

cd $LSFUSION_HOME

case "$1" in
  start)
    /usr/bin/jsvc \
    -cwd $LSFUSION_HOME \
    -home $JAVA_HOME \
    -jvm server \
    -pidfile $PID_FILE \
    -outfile $LSFUSION_HOME/logs/stdout.log \
    -errfile $LSFUSION_HOME/logs/stderr.log \
    $LSFUSION_OPTS \
    -cp $CLASSPATH \
    lsfusion.server.logics.BusinessLogicsBootstrap
    exit $?
    ;;

  stop)
    /usr/bin/jsvc \
    -home $JAVA_HOME \
    -stop \
    -pidfile $PID_FILE \
    lsfusion.server.logics.BusinessLogicsBootstrap
    exit $?
    ;;

  *)
    echo "Usage lsfusion.sh start/stop"
    exit 1;;
esac
